# ansible-pull.service
[Unit]
Description=ansible pull webhook
After=network.target

[Service]
User=http
Type=simple
ExecStart=/bin/webhook -hooks /root/hooks.json -verbose
Restart=always
StartLimitInterval=1min
StartLimitBurst=60

[Install]
WantedBy=multi-user.target
